# deploy/autocierge-mcp.service
[Unit]
Description=Autocierge MCP tool server (localhost only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=autocierge
EnvironmentFile=/etc/autocierge/app.env
ExecStart=/opt/autocierge/mcp-server
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
